# lvl-bypassme/config/supervisor/flask.conf
[program:flask]
command=python3 app.py
directory=/var/www/fileserver
autostart=true
autorestart=true
startretries=3
user=www-data