# lvl-bypassme/config/supervisor/goapi.conf
[program:goapi]
command=/var/www/front/goapi
directory=/var/www/front/
autostart=true
autorestart=true
startretries=3
user=www-data
environment=PROXY_SECRET="DEEZNUTS"